# infra/deploy/launchd/dev.clawdstrike.darwin-telemetry.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.clawdstrike.darwin-telemetry</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/darwin-telemetry-bridge</string>
        <string>--nats-url</string>
        <string>nats://localhost:4222</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>StandardOutPath</key>
    <string>/usr/local/var/log/clawdstrike/darwin-telemetry-bridge.log</string>

    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/clawdstrike/darwin-telemetry-bridge.error.log</string>

    <key>WorkingDirectory</key>
    <string>/usr/local/var/lib/clawdstrike</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>darwin_telemetry_bridge=info</string>
    </dict>

    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>65536</integer>
        <key>NumberOfProcesses</key>
        <integer>4096</integer>
    </dict>

    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>65536</integer>
        <key>NumberOfProcesses</key>
        <integer>4096</integer>
    </dict>

    <key>ProcessType</key>
    <string>Standard</string>

    <key>LowPriorityIO</key>
    <false/>

    <key>Nice</key>
    <integer>0</integer>
</dict>
</plist>
